# tower-content-generation.service
[Unit]
Description=Tower Content Generation API
After=network.target postgresql.service tower-echo-brain.service
Requires=tower-echo-brain.service

[Service]
Type=simple
User=patrick
Group=patrick
WorkingDirectory=/opt/tower-echo-brain
Environment=PYTHONPATH=/opt/tower-echo-brain
ExecStart=/opt/tower-echo-brain/venv/bin/python content_generation_api.py
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# Environment variables
Environment=POSTGRES_HOST=localhost
Environment=POSTGRES_DB=anime_production
Environment=POSTGRES_USER=patrick
Environment=POSTGRES_PASSWORD=RP78eIrW7cI2jYvL5akt1yurE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tower-content-generation

# Resource limits
LimitNOFILE=65536
PrivateTmp=true

[Install]
WantedBy=multi-user.target